# deploy/pi/kic-brand-worker.service
[Unit]
Description=KIC Brand Portal Autonomy Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kic-brand-portal
EnvironmentFile=/opt/kic-brand-portal/.env
ExecStart=/usr/bin/node /opt/kic-brand-portal/server/brand/worker_brand_autonomy.js
Restart=always
RestartSec=30
User=pi

[Install]
WantedBy=multi-user.target
